# overlays/openmptcprouter/common/package/base-files/files/etc/uci-defaults/99-cellular-control-owner
#!/bin/sh

# Keep exactly one userspace owner for each MBIM control port. If a preserved
# configuration selects native netifd MBIM, prevent ModemManager from starting
# during the same boot. If no cellular owner is configured yet, retain this
# uci-defaults script so it can decide after a later configuration restore.

if uci -q show network | grep -q "\\.proto='mbim'$"; then
	/etc/init.d/modemmanager disable >/dev/null 2>&1 || true
	/etc/init.d/modemmanager stop >/dev/null 2>&1 || true
	exit 0
fi

if uci -q show network | grep -q "\\.proto='modemmanager'$"; then
	exit 0
fi

exit 1
